# deploy/com.baidupancopy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.baidupancopy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/baidu-pan-copy/venv/bin/python</string>
        <string>/Users/USERNAME/baidu-pan-copy/start.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/baidu-pan-copy</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/baidu-pan-copy/logs/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/baidu-pan-copy/logs/stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
